# system/frame.service.template
[Unit]
Description=frame

[Service]
User=<USERNAME>
WorkingDirectory=<WORKSPACE>
ExecStart=python3 -u <WORKSPACE>/data/light-frame.py
Restart=always

[Install]
WantedBy=multi-user.target
