# system/player.service.template
[Unit]
Description=player

[Service]
User=<USERNAME>
WorkingDirectory=<WORKSPACE>
ExecStart=player
Restart=always

[Install]
WantedBy=multi-user.target
