# system/startup.service.template
[Unit]
Description=player

[Service]
User=<USERNAME>
WorkingDirectory=<WORKSPACE>
ExecStart=player

[Install]
WantedBy=multi-user.target
